# deploy/entasburada.service
[Unit]
Description=ENTASBURADA Next.js production sunucusu
After=network.target

[Service]
Type=simple
User=entas
WorkingDirectory=/opt/entasburada/apps/web
EnvironmentFile=/opt/entasburada/.env.production.local
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/pnpm start
Restart=always
RestartSec=5
# JSON depolara yazma icin dosya sistemi erisimi gerekli
NoNewPrivileges=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
